# dotfiles/profile.d/bash
#!/bin/bash

#
# To prevent syntastic messages.
#
# Ignore globally until excluding per project via rc file is allowed
# (koalaman/shellcheck#220)
#
export SHELLCHECK_OPTS='-x -e SC1090'
