# dotfiles/profile.d/cabal
#!/bin/bash

pathadd "$HOME/.cabal/bin"
