# dotfiles/profile.d/cabal
#!/bin/bash

path_append "$HOME/.cabal/bin"
