# dotfiles/profile.d/git
#
# Alias hub to git
#
eval "$(hub alias -s)"
