# dotfiles/profile.d/go
#!/bin/bash

path_append "/usr/local/go/bin"
