# dotfiles/profile.d/go
#!/bin/bash

pathadd "/usr/local/go/bin"
