# dotfiles/profile.d/latex
#!/bin/bash

path_append "/usr/local/texlive/2016/bin/x86_64-linux"
infopath_append "/usr/local/texlive/2016/texmf-dist/doc/info"
