# dotfiles/profile.d/latex
#!/bin/bash

PATH_append "/usr/local/texlive/2016/bin/x86_64-linux"
MANPATH_append "/usr/local/texlive/2016/texmf-dist/doc/man"
INFOPATH_append "/usr/local/texlive/2016/texmf-dist/doc/info"
