# dotfiles/profile.d/linuxbrew
#!/bin/bash

path_prepend "$HOME/.linuxbrew/bin"
manpath_prepend "$HOME/.linuxbrew/share/man"
xdg_data_dirs_prepend "$HOME/.linuxbrew/share"
infopath_prepend "$HOME/.linuxbrew/share/info"
