# dotfiles/profile.d/paths
#!/bin/bash

#
# Appends to PATH without duplicates and non folders
#
function path_append {
  append_folder_to_env "$1" "PATH"
}

#
# Appends to MANPATH without duplicates and non folders
#
function manpath_append {
  append_folder_to_env "$1" "MANPATH"

  manpath_keep_configuration_file_list
}

#
# Appends to INFOPATH without duplicates and non folders
#
function infopath_append {
  append_folder_to_env "$1" "INFOPATH"
}

#
# Appends a folder to an environment variable made of a list of folders
# separated by colons
#
function append_folder_to_env {
  if [ -d "$1" ] && [[ ":${!2}:" != *":$1:"* ]]
  then
    export "$2=${!2:+"${!2}:"}$1"
  fi
}

function manpath_keep_configuration_file_list {
  if [[ ! $MANPATH =~ :.* ]]
  then
    export MANPATH=":$MANPATH"
  fi
}
