# dotfiles/profile.d/python
#!/bin/bash

pathadd "$HOME/.local/bin"
