# dotfiles/profile.d/python
#!/bin/bash

path_append "$HOME/.pyenv/bin"

eval "$(pyenv init -)"
eval "$(pyenv virtualenv-init -)"
