# dotfiles/profile.d/python
#!/bin/bash

path_append "$HOME/.local/bin"
