# dotfiles/profile.d/ruby
#
# Enable chruby
#
source /usr/local/share/chruby/chruby.sh

#
# Select default Ruby
#
source /usr/local/share/chruby/auto.sh

#
# JRUBY_OPTS to make jruby startup time acceptable
#
export JRUBY_OPTS='--dev --debug'
