# dotfiles/profile.d/ruby
#
# Enable chruby
#
source $HOME/.linuxbrew/opt/chruby/share/chruby/chruby.sh

#
# Select default Ruby
#
chruby 2.2
