# dotfiles/profile.d/ruby
#
# Enable chruby
#
source /usr/local/share/chruby/chruby.sh

#
# Select default Ruby
#
chruby 2.2
