# dotfiles/profile.d/utils
#!/bin/bash

#
# Recursive search and replace method
#
# @param $1 Old string to be replaced
# @param $2 New string to replace it
#
function search_and_replace {
  [ $# -eq 0 ] && { echo "search_and_replace <old> <new>"; return; }

  git grep -Il '' | xargs perl -pi -e "s/$1/$2/g"
}

#
# Downloads all mp3 files linked to by a web page.
#
# @param $2 Target url
#
# @example download_mp3 http://bohemianbetyars.hu/release/stone-soup/
#
function download_mp3 {
  wget "$2" -O - | grep -oh "http[^ ]*.mp3" | wget -i -
}

#
# Converts a video file to a gif using `imagemagick` + `ffmpeg`.
#
# @param $1 Source video
# @param $2 Target gif
#
function video2gif {
  convert -set delay 10 -layers Optimize "$1" "$2"
}

#
# Adds to PATH without duplicates and non folders
#
function pathadd {
  add_folder_to_env "$1" "PATH"
}

#
# Adds to MANPATH without duplicates and non folders
#
function manpathadd {
  add_folder_to_env "$1" "MANPATH"
}

#
# Adds to INFOPATH without duplicates and non folders
#
function infopathadd {
  add_folder_to_env "$1" "INFOPATH"
}

#
# Adds a folder to an environment variable made of a list of folders separated
# by colons
#
function add_folder_to_env {
  if [ -d "$1" ] && [[ ":${!2}:" != *":$1:"* ]]
  then
    export "$2=${!2:+"${!2}:"}$1"
  fi
}
