# dotfiles/profile.d/utils
#
# Recursive search and replace method
#
# @param $1 Old string to be replaced
# @param $2 New string to replace it
#
function search_and_replace {
  [ $# -eq 0 ] && { echo "search_and_replace <old> <new>"; return; }

  find -type f -name '*.rb' -exec sed -i "s/$1/$2/g" {} \;
}

#
# Downloads all mp3 files linked to by a web page.
#
# @param $2 Target url
#
# @example download_mp3 http://bohemianbetyars.hu/release/stone-soup/
#
function download_mp3 {
  wget "$2" -O - | grep -oh "http[^ ]*.mp3" | wget -i -
}
