# dotfiles/profile.d/vi
#
# Configure alias for vi
#
alias vi=nvim
