# dotfiles/profile.d/vim
#
# Configure alias for vim
#
alias vim=nvim
